# s6-overlay/s6-rc.d/auth-proxy/run
#!/command/with-contenv sh
cd /app/auth-proxy
exec node server.js
