# s6-overlay/s6-rc.d/session-cleanup/run
#!/command/with-contenv sh
# Session cleanup service - kills dtach sessions when no clients connected for > SESSION_TTL seconds

STATE_FILE="/tmp/session-detached-since"
SOCKET_PATH="/home/claude/.dtach/claude-session.sock"
CHECK_INTERVAL=15

# Default to 30 minutes if not set
SESSION_TTL="${CLAUDE_SESSION_TTL:-1800}"

log() {
    echo "[session-cleanup] $*"
}

# Check if any clients are connected via established TCP connections to ttyd
has_ttyd_clients() {
    CONN_COUNT=$(ss -tn state established sport = :${PROXY_PORT:-8080} | tail -n +2 | wc -l)
    CONN_COUNT="${CONN_COUNT:-0}"
    [ "$CONN_COUNT" -gt 0 ]
}

# Check if a dtach session exists
has_dtach_session() {
    test -S "$SOCKET_PATH" && \
    pgrep -f "dtach.*claude-session" >/dev/null 2>&1
}

while true; do
    sleep "$CHECK_INTERVAL"

    # Check if any session exists
    if ! has_dtach_session; then
        rm -f "$STATE_FILE"
        continue
    fi

    # Check if any client is connected via ttyd websocket connections
    if has_ttyd_clients; then
        # Clients connected - reset timer
        if [ -f "$STATE_FILE" ]; then
            log "Client reconnected, resetting timer"
            rm -f "$STATE_FILE"
        fi
    else
        # No clients connected
        NOW=$(date +%s)

        if [ ! -f "$STATE_FILE" ]; then
            # First time seeing no clients, record timestamp
            echo "$NOW" > "$STATE_FILE"
            log "No clients connected, starting TTL countdown (${SESSION_TTL}s)"
        else
            # Check how long since last client
            DETACHED_SINCE=$(cat "$STATE_FILE" 2>/dev/null || echo "$NOW")
            ELAPSED=$((NOW - DETACHED_SINCE))

            if [ "$ELAPSED" -ge "$SESSION_TTL" ]; then
                log "Session exceeded TTL (${ELAPSED}s >= ${SESSION_TTL}s), terminating"

                # Kill dtach server process (the one with ppid=1, i.e. the session host)
                for pid in $(pgrep -f "dtach.*claude-session" 2>/dev/null); do
                    ppid=$(ps -o ppid= -p "$pid" 2>/dev/null | tr -d " ")
                    if [ "$ppid" = "1" ]; then
                        kill -TERM "$pid" 2>/dev/null
                        log "Killed dtach server (PID $pid)"
                    fi
                done

                # Clean up stale socket
                rm -f "$SOCKET_PATH" 2>/dev/null

                rm -f "$STATE_FILE"
            fi
        fi
    fi
done
