# s6-overlay/s6-rc.d/ttyd/run
#!/command/with-contenv sh
exec ttyd -W -p "${PROXY_PORT:-8080}" \
    -t fontSize=14 \
    -t "theme={'background':'#1e1e1e','foreground':'#d4d4d4'}" \
    -t reconnect=3 \
    -t scrollback=10000 \
    -P 30 \
    /home/claude/scripts/claude-session.sh
