# s6-overlay/s6-rc.d/sshd/run
#!/command/with-contenv sh
# Generate host keys if they don't exist
if [ ! -f /etc/ssh/ssh_host_rsa_key ]; then
    ssh-keygen -A
fi
exec /usr/sbin/sshd -D -e
