# s6-overlay/s6-rc.d/init-permissions/up
/bin/sh -c "mkdir -p /home/claude/.claude && if [ $(stat -c%s /home/claude/.claude.json 2>/dev/null || echo 0) -lt 2 ] || ! jq -e 'keys' /home/claude/.claude.json >/dev/null 2>&1; then echo '{}' > /home/claude/.claude.json; fi && chown -R claude:claude /home/claude/.claude /home/claude/workspace /home/claude/.claude.json && SSH_PASS=$(head -c 32 /dev/urandom | base64 | tr -dc 'a-zA-Z0-9' | head -c 24) && echo \"claude:\$SSH_PASS\" | chpasswd && echo \"\$SSH_PASS\" > /run/claude-ssh-pass && chmod 600 /run/claude-ssh-pass && echo 'Claude Code Terminal starting...'"
