# s6-overlay/s6-rc.d/mcp-server/run
#!/command/with-contenv sh

# This server ALWAYS runs: besides the JSON-RPC /mcp endpoint it also serves the
# web auth layer (/login, /auth, /logout) that Caddy forward_auth depends on, so
# gating the whole process on MCP_ENABLED would break login. The MCP_ENABLED flag
# is honored inside server.js, where it gates only the /mcp endpoint + discovery.
cd /app/mcp-server
exec s6-setuidgid claude node server.js
